# vasp-mcp-systemd-services/vasp-query.service
[Unit]
Description=VASP Query MCP Server

[Service]
Type=simple
WorkingDirectory=/home/duguex/vasp_incar
Environment=PATH=/home/duguex/.conda/envs/dgkan_rocm_3.11/bin
ExecStart=/home/duguex/.conda/envs/dgkan_rocm_3.11/bin/python3 /home/duguex/vasp_incar/vasp_query/mcp_server.py --transport http --port 8932 --host 0.0.0.0
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
